# commander-service-files/executor.socket
[Socket]
ListenFIFO=/run/executor.sk

[Install]
WantedBy=sockets.target