# executor-service-files/executor.socket
[Socket]
ListenFIFO=/run/executor.sk
SocketMode=0622

[Install]
WantedBy=sockets.target